# bkp-srv/home/admin/docker/analytics/matomo/build/matomo_archive.cron
MAILTO="webmaster@silene.eu"
5 * * * * www-data /usr/local/bin/php /var/www/html/console core:archive --url=https://analytics.silene.eu/ > /var/www/html/tmp/logs/matomo-archive.log 2> /var/www/html/tmp/logs/matomo-archive-errors.log
