# web-srv/etc/cron.d/grafana
# /etc/cron.d/grafana: backup Grafana configuration

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
# MAILFROM doesn't work now, use alternative below
#MAILTO=adminsys@silene.eu
#MAILFROM=mailer@cbn-alpin.fr

# First day of each month at 00:01AM
1 0 1 * * admin /usr/local/bin/docker-compose run monitor-grafana-backup 2>&1 | mail -s "[Cron] <admin@sinp-paca-web> grafana" -r "mailer@cbn-alpin.fr" adminsys@silene.eu